# rpmbuild/drove-controller-template/etc/init/drove-controller.conf
description "drove-controller"

start on filesystem
stop on runlevel S

respawn
oom never

setuid "drove"
setgid "drove"

exec /usr/bin/drove-controller
